# systemd/cryptodaemon.service.in
[Unit]
Description=Cryptodaemon
Requires=home.mount rm2fb.service
Before=xochitl.service launcher.service sync.service tarnish.service

[Service]
Type=notify
ExecStart=@sbindir@/cryptodaemon -libexecpath @LIBEXECDIR_CRYPTODAEMON@
Restart=on-failure
Environment=LD_PRELOAD=@LIBRM2FB_CLIENT_SO_PATH@
# allow cleanup by new subprocesses after SIGTERM
KillMode=mixed
# do not time out and progress with bootup
TimeoutStartSec=infinity
Environment="PATH=/opt/bin:/opt/sbin:/opt/usr/bin:/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin"

[Install]
WantedBy=multi-user.target
