# dist/etc/systemd/system/cryptodaemon.service
[Unit]
Description=Cryptodaemon
Requires=home.mount rm2fb.service
Before=xochitl.service launcher.service sync.service tarnish.service

[Service]
Type=notify
ExecStart=/home/crypto/bin/cryptodaemon
Restart=on-failure
Environment=LD_PRELOAD=/home/crypto/lib/librm2fb_client.so
# allow cleanup by new subprocesses after SIGTERM
KillMode=mixed
# do not time out and progress with bootup
TimeoutStartSec=infinity

[Install]
WantedBy=multi-user.target
